# coin-service-conf/renec.service.erb
[Unit]
Description=Renec Daemon
After=network.target

[Service]
User=ubuntu
Group=ubuntu
ExecStart=<%= start_command %>
StandardOutput=append:/home/ubuntu/renec-cluster/renec-validator.log
StandardError=append:/home/ubuntu/renec-cluster/renec-validator-error.log
Restart=on-failure

# Specifies which signal to use when killing a service. Defaults to SIGTERM.
KillMode=process
KillSignal=SIGINT
TimeoutStopSec=300
RestartSec=10s

[Install]
WantedBy=default.target
